# crontab
# Wake databases every 6 hours
SHELL=/bin/sh
PATH=/usr/local/bin:/usr/bin:/bin

0 */8 * * * root node /app/wake.js >> /var/log/db-waker.log 2>&1
